# zezere_ignition/zezere_ignition.timer
[Unit]
Description=Trigger Ignition for Zezere until it finishes

[Timer]
OnActiveSec=10sec
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
