# zezere_ignition/zezere_ignition.service
[Unit]
Description=Run Ignition for Zezere

[Service]
Type=oneshot
ExecStart=/usr/bin/zezere_ignition
